# docker/supervisord/supervisord.conf
[unix_http_server]
file=/var/run/supervisor.sock
chmod=0770
chown=nobody:nogroup
# username = mysupervisord
# password = mysupervisordpass

[supervisord]
nodaemon=true
user=root

# Define o comando para iniciar o servidor FrankenPHP
[program:frankenphp]
command = frankenphp php-server -r public/
user = root
autostart = true
